# ansible/roles/quark_server/templates/quark.service.j2
[Unit]
Description=Quark Server
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=3
ExecStart={{ quark_dir }}/target/release/server --metrics prometheus
LimitNOFILE=65536
